# docker/supervisord.conf
; ============================================================================
; supervisord.conf — Sprint 7.3
; ============================================================================
; Avvia nginx + php-fpm (+ queue worker opzionale) nello stesso container.
; Log in stdout/stderr → raccolti da Render/Docker logs.
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---- PHP-FPM ---------------------------------------------------------------
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=QUIT

; ---- nginx -----------------------------------------------------------------
; Usa envsubst per sostituire ${PORT} prima di avviare
[program:nginx]
command=/bin/sh -c "envsubst '$$PORT' < /etc/nginx/nginx.conf > /etc/nginx/nginx.conf.rendered && nginx -c /etc/nginx/nginx.conf.rendered -g 'daemon off;'"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=QUIT

; ---- Queue worker Redis (Sprint 7.6) ---------------------------------------
; Laravel queue worker per job asincroni (email Resend, notifiche, BRT sync).
; In produzione Render: ENABLE_QUEUE_WORKER=false perche' esiste un servizio
; worker dedicato (vedi render.yaml → spedizionefacile-queue).
; In locale / VPS self-hosted: ENABLE_QUEUE_WORKER=true, numprocs=2 per
; parallelismo di base. Alzare a 4+ solo se la coda si allunga.
;
; Flag:
;   --tries=3      retry max 3 volte (Laravel marca failed al 4°)
;   --timeout=60   SIGTERM dopo 60s job
;   --max-jobs=1000 riavvio worker dopo 1000 job (anti memory leak)
;   --max-time=3600 riavvio dopo 1h (pickup deploy)
;   --sleep=3      idle polling 3s (Redis BLPOP tiene viva la connessione)
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --tries=3 --timeout=60 --max-jobs=1000 --max-time=3600 --sleep=3
autostart=%(ENV_ENABLE_QUEUE_WORKER)s
autorestart=true
priority=30
numprocs=2
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopwaitsecs=65  ; timeout job 60s + 5s margine per flush log
